# aim-lead-generator.service
[Unit]
Description=AIM Lead Generator Web Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/lead-generator
ExecStart=/usr/bin/python3 web_app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
